# supervisord.conf
[program:node]
command = node index.js1
directory = /home/dotcloud/current