# supervisord.conf
[program:node]
command = node index.js
directory = /Users/dZ/Projects/node-url-shortener